# guestwlan.target
[Unit]
Description=Guest WLAN Target for QR/Keygen, AP, GUI

[Install]
WantedBy=multi-user.target
